# slack-squid-proxy.service
[Unit]
Description=Slack Squid Proxy Automation Bot
After=network.target squid.service
Wants=squid.service

[Service]
Type=simple
# Run as root: the VM-level no_new_privs flag blocks sudo from escalating
# privileges, making it impossible to reload Squid as a non-root user.
# Filesystem access is still restricted by ProtectSystem and ReadWritePaths.
User=root
WorkingDirectory=/opt/The-Squid-Slack-Bot
EnvironmentFile=/opt/The-Squid-Slack-Bot/.env
ExecStart=/usr/bin/python3 /opt/The-Squid-Slack-Bot/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
# NOTE: NoNewPrivileges is intentionally NOT set — the bot must call
# `sudo squid -k reconfigure` which requires sudo's setuid bit to work.
# Privilege scope is restricted via /etc/sudoers.d/slack-squid-bot instead.
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/etc/squid /etc/squid/lists /etc/squid/conf.d

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
