# slack-squid-proxy.service
[Unit]
Description=Slack Squid Proxy Automation Bot
After=network.target squid.service
Wants=squid.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/The-Squid-Slack-Bot
EnvironmentFile=/opt/The-Squid-Slack-Bot/.env
ExecStart=/usr/bin/python3 /opt/The-Squid-Slack-Bot/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
